# docker-data/crontabs/root
# 每小时统计百度统计地域和访问统计实时数据
0 */1 * * * php /var/www/artisan statistics:bd_tj 2
# 每天12点05分统计统计地域和访问统计昨日数据
5 0 * * * php /var/www/artisan statistics:bd_tj 1
# 每月第一天创建操作日志表
0 0 1 * * php /var/www/artisan createOperationLogTable
