# docker-data/crontabs/root
# 下个月十分钟前创建操作日志表
50 23 28-31 * * php /var/www/artisan autoCreateTable >> /var/www/storage/logs/autoCreateTable.log
