# docker-data/crontabs/root
# DO NOT EDIT THIS FILE - edit the master and reinstall.
0 1 * * * php /var/www/artisan statistics:bd_tj 2
0 */60 * * * php /var/www/artisan statistics:bd_tj 1
0 0 1 * * php /var/www/artisan createOperationLogTable
