# docker-data/crontabs/root
# 月底最后一分钟创建操作日志表
59 23 28-31 * * php /var/www/artisan autoCreateTable >> /var/www/storage/logs/autoCreateTable.log
